# orquestador_sistema/etc/systemd/system/orchestrator.service.d/wg.conf
[Service]
# Fix 2026-05-25: el WG automático del cliente (vpn.py) escala vía `sudo` al
# helper /usr/local/sbin/orchestrator-wg-helper. NoNewPrivileges=true bloqueaba
# sudo ("no new privileges flag is set") → el túnel nunca se levantaba en el
# preload. Lo desactivamos para este servicio (el sudoers NOPASSWD está acotado
# solo al helper, que solo gestiona interfaces eng_*, así que el blast radius
# es limitado).
NoNewPrivileges=false
# ProtectSystem=strict deja /etc read-only en el namespace del servicio;
# permitir que el helper (root) escriba el .conf del túnel.
ReadWritePaths=/etc/wireguard
