# orquestador_sistema/etc/cron.d/audit-backup
0 3 * * * root /usr/local/sbin/audit-backup.sh
